# services/ldap.lab.lan/context/s6/php/run
#!/bin/sh

cd /www
# Use `su-exec $USER` to change permissions.
/usr/sbin/php-fpm82 -FO